# .config/autostart/keyboards.sh.desktop
[Desktop Entry]
Exec=/home/julio/bin/keyboards.sh
Icon=dialog-scripts
Name=keyboards.sh
Path=
Type=Application
X-KDE-AutostartScript=true
